# backend/com.theofficecms.flask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.theofficecms.flask</string>
	<key>Program</key>
	<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/venv/bin/python</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/venv/bin/python</string>
		<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/app.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend</string>
	<key>StandardOutPath</key>
	<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/logs/flask.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/logs/flask-error.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/monicanieckula/Documents/GitHub/theOfficeCMS/backend/venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Restart</key>
	<string>on-demand</string>
</dict>
</plist>
